# private_dot_config/systemd/user/xdg-desktop-portal-hyprland.service.d/override.conf
[Service]
ExecStartPre=/bin/sleep 2
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=5
